# upstart/elastic-search.conf
description "Elastic Search"
start on starting tripolist
stop on stopping tripolist
respawn
exec sudo /etc/init.d/elasticsearch start