# upstart/tripolist-upstart.conf
description "Tripolist application"
start on runlevel [2345]
stop on runlevel [016]

pre-start script
  mkdir -p /var/log/tripolist
  chown -R www-data /var/log/tripolist

  mkdir -p /var/run/tripolist
  chown -R www-data /var/run/tripolist
end script


